# deploy/templates/gym-bot.service.j2
[Unit]
Description=Gym Bot
After=network.target docker.service
Requires=docker.service

[Service]
User={{ deploy_user }}
WorkingDirectory={{ project_dir }}
EnvironmentFile={{ project_dir }}/.env
ExecStart={{ project_dir }}/.venv/bin/python -m gym_bot
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
